# 自动化/Saltstack/资料/Deploy_Supervisor/supervisor/templates/program.conf.tmpl
;;; File managed by Salt
[program:{{ program }}]
{% for key, value in values.items() -%}
{{ key }}={{ value }}
{% endfor -%}

